# systemd/vpskit-worker.service
[Unit]
Description=VPSKit Provisioning Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vpskit
Group=vpskit
WorkingDirectory=/opt/vpskit
EnvironmentFile=/etc/vpskit.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/vpskit/venv/bin/python -m vpskit.provisioning
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vpskit-worker
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
StateDirectory=vpskit

[Install]
WantedBy=multi-user.target
